# packaging/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop and disable the service before removal
        if systemctl is-active --quiet cortex 2>/dev/null; then
            systemctl stop cortex || true
        fi
        if systemctl is-enabled --quiet cortex 2>/dev/null; then
            systemctl disable cortex || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
